# test_meter/supervisord.conf
[program:dbmetricsaggregator]
command={{ PYTHON }} manage.py dbmetermetricsaggregator
autoreload=True
